# conf/supervisor.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/mysql-entrypoint.sh mysqld
stdout_logfile=/logs/%(program_name)s.log
stderr_logfile=/logs/%(program_name)s.log

[program:cron]
command=cron -f -L 15
autostart=true
autorestart=true
startretries=0
stdout_logfile=/logs/%(program_name)s.log
stderr_logfile=/logs/%(program_name)s.log